# net-vpn/softether/files/softether-vpnserver.initd
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
}

start() {
	ebegin "Starting SoftEther VPN Server"
	/opt/softether/vpnserver/vpnserver start >/dev/null 2>&1
	eend $?
}

stop() {
	ebegin "Stopping SoftEther VPN Server"
	/opt/softether/vpnserver/vpnserver stop >/dev/null 2>&1
	eend $?
}
