# net-vpn/softether/files/softether-vpnclient.initd
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
}

start() {
	ebegin "Starting SoftEther VPN client"
	/opt/softether/vpnclient/vpnclient start
	eend $?
}

stop() {
	ebegin "Stopping SoftEther VPN client"
	/opt/softether/vpnclient/vpnclient stop
	eend $?
}
