# net-proxy/mtproxy/files/mtproxy-config.service
[Unit]
Description=Update MTProxy configuration

[Service]
Type=oneshot
ExecStart=/bin/systemctl try-restart mtproxy.service
